# files/etc/uci-defaults/90_setup.sh
cd /root 
git clone https://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@github.com/metamericdesign/gsSystems_Lighting.git --depth 1
cd /srv/www
git clone https://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@github.com/metamericdesign/BaseStationPHP.git --depth 1

chmod 777 /srv/www/BaseStationPHP/Php/LightingNodeLog.php

exit 0