# files/etc/init.d/runLightWorker
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=97

start_service() {
    procd_open_instance
    procd_set_param command /usr/bin/python3 /root/gsSystems_Lighting/win_light.py
    procd_close_instance
}